# samples/vorimport.supervisor.conf
[program:vorimport]
command=/opt/vorimport/current/bin/vorimport
user=root
autostart=true
autorestart=true