# samples/vorimport.supervisor.conf
[program:vorimport]
command=/opt/vorimport/cirrent/vorimport
user=root
autostart=true
autorestart=true