# samples/vorimport.supervisor.conf
[program:vorimport]
command=/opt/vorimport/current/vorimport
directory=/opt/vorimport/current/
user=root
autostart=true
autorestart=true